# config/systemd/user/openrgb.service
[Unit]
Description=OpenRGB Server

[Service]
Type=simple
ExecStart=/usr/bin/openrgb --server

[Install]
WantedBy=default.target

